# deploy/agenc-mcp-shared-memory.service
[Unit]
Description=AgenC MCP shared-memory scratchpad (via router)
After=network.target mcp-router.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/AgenC/integrations
ExecStart=/root/.local/bin/uv run agenc-shared-memory-mcp --host 127.0.0.1 --port 9102 --router http://127.0.0.1:9003
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
