# deploy/agenc-mcp-web-search.service
[Unit]
Description=AgenC MCP web-search (DuckDuckGo via router)
After=network.target
Requires=mcp-router.service
After=mcp-router.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/AgenC/integrations
# Kill stray listeners on this port (manual MCP runs / zombie PIDs) so bind does not fail with EADDRINUSE
ExecStartPre=-/bin/sh -c 'fuser -k 9101/tcp 2>/dev/null || true; sleep 0.5'
ExecStart=/root/.local/bin/uv run agenc-web-search-mcp --host 127.0.0.1 --port 9101 --router http://127.0.0.1:9003
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
