# deploy/worker3.service
[Unit]
Description=AgenC Worker 3 — Sentiment Analyst
After=network.target agenc-api.service axl-emitter.service axl-worker3.service mcp-router.service agenc-mcp-web-search.service agenc-mcp-shared-memory.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/AgenC
EnvironmentFile=/root/AgenC/.env
ExecStart=/root/.local/bin/uv run python worker3.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
